# systemd/auto-patch.service
[Unit]
Description=Install available patches

[Service]
Type=oneshot
ExecStart=/usr/sbin/auto-patch
